# dev/wireplumber/evo8/evo8-setup.service
[Unit]
Description=EVO8 Audio Setup
After=pipewire.service wireplumber.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/evo8-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
